# roles/supervisor/templates/laravel.conf
[program:laravel]
command=php artisan queue:listen
directory=/var/www/laravel/
stdout_logfile=/var/www/laravel/app/storage/logs/supervisord.log
redirect_stderr=true
